# s6/services.d/telly/run
#!/usr/bin/with-contenv bash

source /etc/env

exec s6-setuidgid telly /usr/bin/telly
